# ai-photo.service
[Unit]
Description=AI ID Photo Generation Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/photo-service
ExecStart=/root/photo-service/venv/bin/python3 /root/photo-service/ai_photo.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
